# modules/hrt/templates/hrt-worker.conf
start on starting hrt
stop on stopping hrt
